# deploy/cron.d/alpaca-bot
SHELL=/bin/bash
# Commands are scheduled at both possible UTC hours and gated on New York wall time
# by run_if_ny_time.sh. This keeps paper proof checks aligned through DST changes.
# Paper readiness — primary run before NYSE open, retry before the 10:00 ET entry window,
# then short-grace repair passes around the open so a supervisor restart after 09:58
# does not leave the proof gate blocked until the later activity check.
20 13,14 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 0920 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_check.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
55 13,14 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 0955 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_check.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
58 13,14 * * 1-5 root RUN_IF_NY_TIME_GRACE_MINUTES=1 /workspace/alpaca_bot/scripts/run_if_ny_time.sh 0958 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
2 14,15 * * 1-5 root RUN_IF_NY_TIME_GRACE_MINUTES=1 /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1002 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
# Post-open paper activity — diagnostic first pass after the first completed
# 15-minute entry bar has time to settle, fatal second pass mid-session.
25 14,15 * * 1-5 root PAPER_ACTIVITY_CLOSE_ONLY_ON_FAILURE=false /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1025 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_activity /var/lock/alpaca-bot-paper-activity.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_activity_check.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-activity.log 2>&1
0 16,17 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1200 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_activity /var/lock/alpaca-bot-paper-activity.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_activity_check.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-activity.log 2>&1
# Mid-session and post-close readiness refreshes keep the 180-minute
# proof-readiness freshness guard valid through the entry window, post-close
# session guard, profit probe, and proof-status due time. The post-close pass
# skips only the session-unblocked gate because the regular-session flatten is
# expected to have disabled entries for that same proof session.
15 16,17 * * 1-5 root PAPER_READINESS_FORCE_REFRESH=true /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1215 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
25 18,19 * * 1-5 root PAPER_READINESS_FORCE_REFRESH=true /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1425 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
55 20,21 * * 1-5 root PAPER_READINESS_FORCE_REFRESH=true PAPER_READINESS_REQUIRE_SESSION_UNBLOCKED=false /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1655 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
# Post-close proof checks and nightly evolve pipeline.
# apply_candidate.sh chains after nightly (&&) and restarts the supervisor if params changed.
# The env file must be sourced first: deploy/compose.yaml interpolates ${TRADING_MODE} etc.
# from the calling shell, exactly like scripts/deploy.sh does.
10 21,22 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1710 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh session_guard /var/lock/alpaca-bot-session-guard.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/session_guard.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-session-guard.log 2>&1
20 21,22 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1720 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_profit_probe /var/lock/alpaca-bot-profit-probe.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_profit_probe.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-profit-probe.log 2>&1
24 21,22 * * 1-5 root PAPER_READINESS_FORCE_REFRESH=true PAPER_READINESS_REQUIRE_SESSION_UNBLOCKED=false /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1724 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_readiness /var/lock/alpaca-bot-paper-readiness.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_readiness_if_needed.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-paper-readiness.log 2>&1
28 21,22 * * 1-5 root PROOF_STATUS_FAIL_ON_ISSUES=true /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1728 /workspace/alpaca_bot/scripts/run_locked_check_with_audit.sh paper_proof_status /var/lock/alpaca-bot-proof-status.lock /etc/alpaca_bot/alpaca-bot.env /workspace/alpaca_bot/scripts/paper_proof_status.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-proof-status.log 2>&1
30 21,22 * * 1-5 root /workspace/alpaca_bot/scripts/run_if_ny_time.sh 1730 flock -n /var/lock/alpaca-bot-nightly.lock bash -lc 'mkdir -p /var/lib/alpaca-bot/nightly && cd /workspace/alpaca_bot && set -a && source /etc/alpaca_bot/alpaca-bot.env && set +a && docker compose -f deploy/compose.yaml run --rm nightly >> /var/log/alpaca-bot-nightly.log 2>&1 && ./scripts/apply_candidate.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-nightly.log 2>&1'
