# deploy/cron.d/alpaca-bot
# Nightly evolve pipeline — runs 30 min after NYSE close (22:30 UTC = 5:30 PM ET year-round)
# apply_candidate.sh chains after nightly (&&) and restarts the supervisor if params changed.
30 22 * * 1-5 root mkdir -p /var/lib/alpaca-bot/nightly && cd /workspace/alpaca_bot && docker compose -f deploy/compose.yaml run --rm nightly >> /var/log/alpaca-bot-nightly.log 2>&1 && ./scripts/apply_candidate.sh /etc/alpaca_bot/alpaca-bot.env >> /var/log/alpaca-bot-nightly.log 2>&1
